# supervisor.conf
[program:httpstatus_history]

command = php console.php /controllers/internals/history check_history
directory = /var/www/html/httpstatus/
autostart = true
autorestart = true